# relay_supervisor.conf
[program:relay]
user=root
command=npm run start-ecs --prefix /relay/
startretries=999999999999999999999999999
autostart=true
autorestart=true